# agent/rapisys-agent.service
[Unit]
Description=RaPiSys host agent (allowlisted privileged operations)
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/rapisys/agent.env
ExecStart=/usr/bin/node /opt/rapisys/agent/rapisys-agent.cjs
Restart=on-failure
RestartSec=3
# Hardening (the agent needs root for mounts/apt/fan, but we still confine it)
NoNewPrivileges=no
# ProtectHome is disabled because the SunFounder Pironman installer resolves
# HOME from root's passwd entry (/root) and clones into $HOME/pironman5; on this
# platform neither ProtectHome=read-only + ReadWritePaths=/root nor
# ProtectHome=tmpfs yields a writable /root (both leave it read-only — verified
# on-device), so the clone fails. ProtectHome=no is the only setting that lets
# it write. This is an acceptable trade: the agent already runs as root for
# apt/systemd/fan/mount work and reads no config from home dirs (its env comes
# from /etc/rapisys/agent.env), so home sandboxing added little here.
ProtectHome=no
PrivateTmp=yes
RuntimeDirectory=rapisys
RuntimeDirectoryMode=0755
# Keep /run/rapisys across restarts: removing it would orphan the docker
# bind mount and silently sever the container <-> agent socket.
RuntimeDirectoryPreserve=yes

[Install]
WantedBy=multi-user.target
