# agent/rapisys-agent.service
[Unit]
Description=RaPiSys host agent (allowlisted privileged operations)
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/rapisys/agent.env
ExecStart=/usr/bin/node /opt/rapisys/agent/rapisys-agent.cjs
Restart=on-failure
RestartSec=3
# Hardening (the agent needs root for mounts/apt/fan, but we still confine it)
NoNewPrivileges=no
ProtectHome=read-only
# The SunFounder Pironman installer resolves HOME from root's passwd entry
# (/root) and clones into $HOME/pironman5; ProtectHome=read-only would block
# that. Grant write access to /root specifically so the installer's clone
# step works without relaxing home protection elsewhere.
ReadWritePaths=/root
PrivateTmp=yes
RuntimeDirectory=rapisys
RuntimeDirectoryMode=0755
# Keep /run/rapisys across restarts: removing it would orphan the docker
# bind mount and silently sever the container <-> agent socket.
RuntimeDirectoryPreserve=yes

[Install]
WantedBy=multi-user.target
