# agent/rapisys-agent.service
[Unit]
Description=RaPiSys host agent (allowlisted privileged operations)
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/rapisys/agent.env
ExecStart=/usr/bin/node /opt/rapisys/agent/rapisys-agent.cjs
Restart=on-failure
RestartSec=3
# Hardening (the agent needs root for mounts/apt/fan, but we still confine it)
NoNewPrivileges=no
# ProtectHome=tmpfs mounts an empty *writable* tmpfs over /home and /root so the
# real home dirs stay hidden, while the SunFounder Pironman installer — which
# resolves HOME from root's passwd entry (/root) and clones into
# $HOME/pironman5 — can still write its throwaway build clone. (ReadWritePaths
# does NOT re-grant write under ProtectHome=read-only, so tmpfs is the fix.)
# The agent reads no config from home dirs (its env is /etc/rapisys/agent.env).
ProtectHome=tmpfs
PrivateTmp=yes
RuntimeDirectory=rapisys
RuntimeDirectoryMode=0755
# Keep /run/rapisys across restarts: removing it would orphan the docker
# bind mount and silently sever the container <-> agent socket.
RuntimeDirectoryPreserve=yes

[Install]
WantedBy=multi-user.target
